# conf/supervisor/conf.d/magento-consumers.conf
# conf/supervisor/conf.d/magento-consumers.conf
# incloure a /etc/supervisor/conf.d/ al Dockerfile

[supervisord]
nodaemon=true
logfile=/var/www/html/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/www/html/var/log/supervisor


# [include]
# files = /etc/supervisor/conf.d/*.conf


[group:magento-consumers]
programs=magento-consumer-async-operations-all,
    # magento-consumer-nacento-gallery-bulk-async,
    magento-consumer-inventory-reservations-cleanup,
    magento-consumer-inventory-mass-update,
    magento-consumer-inventory-indexer-sourceItem,
    magento-consumer-inventory-reservations-update,
    magento-consumer-inventory-reservations-updateSalabilityStatus,
    magento-consumer-inventory-indexer-stock,
    magento-consumer-inventory-source-items-cleanup,
    magento-consumer-media-content-synchronization,
    magento-consumer-media-gallery-renditions-update,
    # magento-consumer-save-config,
    magento-consumer-media-gallery-synchronization,
    magento-consumer-media-storage-catalog-image-resize,
    magento-consumer-catalog-website-attribute-value-sync,
    magento-consumer-product-action-attribute-update,
    magento-consumer-product-action-attribute-website-update,
    # magento-consumer-export-processor,
    # magento-consumer-sales-rule-quote-trigger-recollect,
    # magento-consumer-sales-rule-update-coupon-usage,
    # magento-consumer-product-alert,
    magento-consumer-codegenerator-processor
priority=10

; ========= ASYNC / BULK =========
[program:magento-consumer-async-operations-all]
command=php bin/magento queue:consumers:start async.operations.all
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_async_operations_all.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_async_operations_all.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

# ; el consumer es "nacento.gallery.process.batch.consumer" (la cua es nacento.gallery.process.batch.queue)
# [program:magento-consumer-nacento-gallery-bulk-async]
# command=/usr/local/bin/php bin/magento queue:consumers:start nacento.media-gallery.sync.sku
# directory=/var/www/html
# user=app
# autostart=true
# autorestart=true
# startsecs=3
# stopsignal=INT
# stopwaitsecs=20
# stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_nacento_gallery_bulk_async.log
# stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_nacento_gallery_bulk_async.err
# stdout_logfile_maxbytes=50MB
# stderr_logfile_maxbytes=50MB
# stdout_logfile_backups=5
# stderr_logfile_backups=5
# killasgroup=true
# stopasgroup=true
# startretries=999

; ========= INVENTORY =========
[program:magento-consumer-inventory-indexer-sourceItem]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.indexer.sourceItem
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_indexer_sourceItem.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_indexer_sourceItem.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-indexer-stock]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.indexer.stock
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_indexer_stock.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_indexer_stock.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-mass-update]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.mass.update
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_mass_update.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_mass_update.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-reservations-cleanup]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.reservations.cleanup
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_cleanup.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_cleanup.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-reservations-update]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.reservations.update
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_update.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_update.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-reservations-updateSalabilityStatus]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.reservations.updateSalabilityStatus
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_updateSalabilityStatus.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_reservations_updateSalabilityStatus.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-inventory-source-items-cleanup]
command=/usr/local/bin/php bin/magento queue:consumers:start inventory.source.items.cleanup
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_source_items_cleanup.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_inventory_source_items_cleanup.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= MEDIA / CONTENT =========
[program:magento-consumer-media-content-synchronization]
command=/usr/local/bin/php bin/magento queue:consumers:start media.content.synchronization
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_media_content_synchronization.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_media_content_synchronization.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-media-gallery-renditions-update]
command=/usr/local/bin/php bin/magento queue:consumers:start media.gallery.renditions.update
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_media_gallery_renditions_update.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_media_gallery_renditions_update.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-media-gallery-synchronization]
command=/usr/local/bin/php bin/magento queue:consumers:start media.gallery.synchronization
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_media_gallery_synchronization.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_media_gallery_synchronization.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-media-storage-catalog-image-resize]
command=/usr/local/bin/php bin/magento queue:consumers:start media.storage.catalog.image.resize
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_media_storage_catalog_image_resize.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_media_storage_catalog_image_resize.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= CATALOG =========
[program:magento-consumer-catalog-website-attribute-value-sync]
command=/usr/local/bin/php bin/magento queue:consumers:start catalog_website_attribute_value_sync
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_catalog_website_attribute_value_sync.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_catalog_website_attribute_value_sync.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-product-action-attribute-update]
command=/usr/local/bin/php bin/magento queue:consumers:start product_action_attribute.update
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_product_action_attribute_update.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_product_action_attribute_update.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-product-action-attribute-website-update]
command=/usr/local/bin/php bin/magento queue:consumers:start product_action_attribute.website.update
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_product_action_attribute_website_update.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_product_action_attribute_website_update.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= EXPORT / CODEGEN =========
; noms correctes: exportProcessor / codegeneratorProcessor
[program:magento-consumer-export-processor]
command=/usr/local/bin/php bin/magento queue:consumers:start exportProcessor
directory=/var/www/html
user=app
autostart=false
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_exportProcessor.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_exportProcessor.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-codegenerator-processor]
command=/usr/local/bin/php bin/magento queue:consumers:start codegeneratorProcessor
directory=/var/www/html
user=app
autostart=false
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_codegeneratorProcessor.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_codegeneratorProcessor.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= SALES / RULES =========
[program:magento-consumer-sales-rule-quote-trigger-recollect]
command=/usr/local/bin/php bin/magento queue:consumers:start sales.rule.quote.trigger.recollect
directory=/var/www/html
user=app
autostart=false
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_sales_rule_quote_trigger_recollect.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_sales_rule_quote_trigger_recollect.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

[program:magento-consumer-sales-rule-update-coupon-usage]
command=/usr/local/bin/php bin/magento queue:consumers:start sales.rule.update.coupon.usage
directory=/var/www/html
user=app
autostart=false
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_sales_rule_update_coupon_usage.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_sales_rule_update_coupon_usage.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= PRODUCT ALERT =========
[program:magento-consumer-product-alert]
command=/usr/local/bin/php bin/magento queue:consumers:start product_alert
directory=/var/www/html
user=app
autostart=false
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_product_alert.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_product_alert.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999

; ========= SAVE CONFIG =========
[program:magento-consumer-save-config]
command=/usr/local/bin/php bin/magento queue:consumers:start saveConfigProcessor
directory=/var/www/html
user=app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopwaitsecs=20
stdout_logfile=/var/www/html/var/log/magento-consumers/consumer_saveConfigProcessor.log
stderr_logfile=/var/www/html/var/log/magento-consumers/consumer_saveConfigProcessor.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
killasgroup=true
stopasgroup=true
startretries=999
